# ansible/roles/agents/templates/kimi-web.service.j2
[Unit]
Description=Kimi Code web UI
Documentation=https://moonshotai.github.io/kimi-code/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=KIMI_CODE_HOME=%h/.kimi-code
Environment=HOME=%h
ExecStart=%h/.local/bin/kimi web --foreground --no-open --port 58628
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
